# health-dashboard.service
[Unit]
Description=Hermes Health Dashboard (Streamlit)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/garmin-health/dashboard
ExecStart=/usr/local/lib/hermes-agent/venv/bin/streamlit run dashboard.py --server.port 8093 --server.headless true
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
